# systemd/rpmsg-service.service
[Unit]
Description=RPMsg Bridge Service (Linux <-> M4F)
Documentation=https://github.com/damianbubak92/bramka-setup
After=network.target
# Czekaj na M4F remoteproc zanim wystartujesz - jeśli M4F nie running,
# OpenTransport nie znajdzie /dev/rpmsg* i service padnie od razu.
# (cold-boot race obsłużony w transport.go: waitForM4FChrdev czeka do 20s)

# Safety net - jeśli serwis pada >3 razy w 60s, systemd przestaje restartować.
# MUSI być w [Unit] (w [Service] nowoczesny systemd to ignoruje!).
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=notify
ExecStart=/opt/bramka/rpmsg-service/rpmsg-service -test hello
WorkingDirectory=/opt/bramka/rpmsg-service

# === NON-ROOT (least privilege) ===
# Patrz modules/08-hardening.sh: user 'bramka' + udev (/dev/rpmsg* -> grupa bramka)
# + path-unit dla czystego reboota (serwis nie ma roota ani CAP_SYS_BOOT).
User=bramka
Group=bramka

# === RESTART POLICY ===
# on-failure (NIE always!) - restartuj tylko gdy crash/non-zero exit.
Restart=on-failure
RestartSec=2

# === WATCHDOG ===
# systemd zabije serwis SIGABRT jeśli nie wyśle WATCHDOG=1 w ciągu 10s.
# Go service kickuje co 5s (WATCHDOG_USEC/2). Działa też jako non-root (sd_notify).
WatchdogSec=10s

# === TIMEOUTS ===
TimeoutStartSec=30s
TimeoutStopSec=15s

# === LOGGING ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpmsg-service

# === SECURITY HARDENING ===
# Zero capabilities - serwis nic nie potrzebuje (device przez grupę, reboot przez
# path-unit, stan przez ReadWritePaths). NIE ustawiać PrivateDevices - musi widzieć
# /dev/rpmsg. NIE ustawiać MemoryDenyWriteExecute - może wadzić Go/cgo.
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/bramka /run/bramka
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
